# competence-installer/Linux/Ubuntu_Debian/db-competence-installer-0.1/etc/systemd/system/db-competence.service
[Unit]
Description=Daemon of Competence Database
After=syslog.target

[Service]
Type=simple
User=(USER)
Group=(GROUP)
WorkingDirectory=(WORKINGDIR)
ExecStart=(WORKINGDIR)run.py
StandardOutput=syslog
StandardError=syslog

[Install]
WantedBy=multi-user.target
